# theme/luci-theme-xargon/root/etc/uci-defaults/luci-theme-xargon
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Xargon=/luci-static/xargon
	set luci.main.mediaurlbase=/luci-static/xargon
	commit luci
EOF
exit 0
