# theme/luci-theme-darkmatter/root/etc/uci-defaults/luci-theme-darkmatter
#!/bin/sh
uci batch <<-EOF
	set luci.themes.darkmatter=/luci-static/darkmatter
	set luci.main.mediaurlbase=/luci-static/darkmatter
	commit luci
EOF
